# config/superlocal-engine.service
[Unit]
Description=Superlocal Memory C++ Backend Engine
Before=superlocal-fuse.service

[Service]
Type=simple
ExecStart=%h/.local/bin/memory_engine --db-path=%h/.superlocal/memory.db --shm-name=superlocal_shm
Restart=on-failure
LimitMEMLOCK=infinity
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=default.target
